# deploy/supervisor-eo-sanasini.conf
# EO Sanasini — Supervisor config (production build)
# Pakai Nuxt production build (bukan dev) untuk hemat resource
# Reload setelah ubah: supervisorctl reread && supervisorctl update

[program:eo-sanasini-laravel]
command=php artisan serve --host=0.0.0.0 --port=8001
directory=/var/www/eo-sanasini/backend
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/eo-sanasini-laravel.log
stderr_logfile=/var/log/supervisor/eo-sanasini-laravel-err.log
stopwaitsecs=10
environment=HOME="/root"

[program:eo-sanasini-nuxt]
command=node .output/server/index.mjs
directory=/var/www/eo-sanasini/frontend
user=root
environment=HOME="/root",NODE_ENV="production",NUXT_PUBLIC_API_BASE="http://194.233.90.53:8001/api",HOST="0.0.0.0",PORT="3001"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/eo-sanasini-nuxt.log
stderr_logfile=/var/log/supervisor/eo-sanasini-nuxt-err.log
stopwaitsecs=10
